# docker/supervisor/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info
pidfile=/run/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

# ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

# ── PHP-FPM ───────────────────────────────────────────────────────────────────
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/php-fpm.log
stderr_logfile=/var/log/supervisor/php-fpm-error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB

# ── Queue Worker: default ─────────────────────────────────────────────────────
[program:queue-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis \
    --queue=default \
    --sleep=3 \
    --tries=3 \
    --max-time=3600 \
    --max-jobs=1000 \
    --memory=256
autostart=true
autorestart=true
numprocs=2
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=60
killasgroup=true

# ── Queue Worker: notifications ───────────────────────────────────────────────
[program:queue-notifications]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis \
    --queue=notifications \
    --sleep=3 \
    --tries=5 \
    --backoff=10 \
    --max-time=3600 \
    --max-jobs=500 \
    --memory=128
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-notifications.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=30
killasgroup=true

# ── Queue Worker: cache (WarmProjectCache) ────────────────────────────────────
[program:queue-cache]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work redis \
    --queue=cache \
    --sleep=5 \
    --tries=2 \
    --max-time=3600 \
    --max-jobs=200 \
    --memory=512
autostart=true
autorestart=true
numprocs=1
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/queue-cache.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=120
killasgroup=true

# ── Laravel Reverb (WebSocket server) ────────────────────────────────────────
[program:reverb]
command=php /var/www/html/artisan reverb:start \
    --host=0.0.0.0 \
    --port=8080 \
    --no-interaction
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=10

[group:workers]
programs=queue-default,queue-notifications,queue-cache
